# deploy/divd-boot-hello.service
[Unit]
Description=DIVD: print a hello-world receipt at boot
After=divd-print.service systemd-udev-settle.service
Wants=divd-print.service

[Service]
Type=oneshot
User=divd
Group=divd
SupplementaryGroups=lp
WorkingDirectory=/opt/divdprint
ExecStart=/opt/venv/bin/python manage.py print_hello
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
